# docker/php/supervisor/conf.d/asyncETL.conf
[program:asyncETL]
command=php /srv/app/bin/console messenger:consume asyncETL
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/asyncETL_stderr.log
